# inst/root_overlay/etc/cron.d/archiveJobs
# Archive old jobs.
# Run Saturday at 1:15 PM GMT
#m h  dom mon dow user cmd
15 13  *   *   6  sg   /sgm/bin/archiveJobs.R
