# inst/root_overlay/etc/cron.d/archiveJobs
# Archive old jobs.
# Run daily at 1:15 PM GMT
#m h  dom mon dow user cmd
15 13  *   *   *  sg   /sgm/bin/archiveJobs.R
