# inst/root_overlay/etc/cron.d/monthly_trash_dump
# Permanently delete all trash folders older than 31 days.
# Run at 4:01 AM GMT on the 1st day of each month.
#m h    dom mon dow    user   cmd
1  4     1   *   *     sg    [ -x /sgm_local/bin/monthly_trash_dump.sh ] && /sgm_local/bin/monthly_trash_dump.sh
