# inst/root_overlay/etc/cron.d/refreshMotusMetaDB
# Run daily at 8:30 AM GMT.
#m h    dom mon dow    user   cmd
30 *     *   *   *     sg     /sgm/bin/refreshMotusMetaDB.R
